# etc/supervisor/conf.d/apache.conf

[program:apache2]
command=/usr/bin/pidproxy /var/run/httpd.pid /bin/bash -c "/usr/sbin/httpd -DFOREGROUND -k start"
redirect_stderr=true
user=root
autorestart=true
priority=400
