# etc/systemd/system/rivendell.service.d/override.conf
[Service]
Environment=DISPLAY=:19.0
ExecStart=
ExecStart=/bin/sh -lc /usr/sbin/rdservice
User=rd
